# data/indicators.target
[Unit]
Description=Target representing the lifecycle of the indicators. Each indicator should be bound to it in its individual service file.

[Install]
# Old name for this target, kept for compatibility
Alias=indicators-pre.target
