# data/indicators-pre.target
[Unit]
Description=All indicator hosts should bind to this target to insure that the indicators will start after they've already started
